# supervisord.conf
[program:uwsgi]
command = uwsgi --ini /home/docker/code/uwsgi.ini
autostart = true
autorestart = true
redirect_stderr = true
stopsignal = QUIT

[supervisord]
logfile = /logs/supervisord.log
logfile_maxbytes = 10MB
logfile_backups = 5
loglevel = info
pidfile = run/supervisord.pid

[supervisorctl]
serverurl = unix://run/supervisor.sock

[unix_http_server]
file = run/supervisor.sock
chmod = 0777

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command = nginx -g "daemon off;"
autostart = true
autorestart = true
redirect_stderr = true
stopsignal = QUIT